# scripts/com.afs.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.afs.supervisor</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace with your AFS venv python path -->
        <string>$HOME/src/lab/afs/.venv/bin/python</string>
        <string>-m</string>
        <string>afs.agents.supervisor</string>
        <string>--max-runs</string>
        <string>1</string>
        <string>--quiet</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Run every 15 minutes -->
    <key>StartInterval</key>
    <integer>900</integer>

    <!-- Also run on load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/afs-supervisor.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/afs-supervisor.stderr.log</string>

    <!-- Resource limits — keep it lightweight -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- Don't restart too fast on failure -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- Keep alive only if launched -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
